# templates/upstart.conf
description "{{name}}"
author "Chris Holcombe and Chris MacNaughton"

start on runlevel [2345]
stop on runlevel [016]

respawn

exec {{vault_path}} server {{vault_options}}

